# install/boesch_backend.service
[Unit]
Description=Morioo MMS — Backend Python FastAPI
After=network.target time-sync.target
Wants=time-sync.target

[Service]
User=ode
Group=ode
WorkingDirectory=/home/ode/boesch_os
Environment="HOME=/home/ode"
Environment="PYTHONPATH=/home/ode/boesch_os"
ExecStart=/home/ode/boesch_os/venv/bin/python main.py
Restart=always
RestartSec=3
# Type=notify : systemd attend le READY=1 avant de considérer le service
# démarré, et active NOTIFY_SOCKET pour que les pings WATCHDOG=1 arrivent.
# Sans cette ligne, WatchdogSec tue le process après 30s (NOTIFY_SOCKET absent).
Type=notify
NotifyAccess=main
WatchdogSec=30

[Install]
WantedBy=multi-user.target
